# podman/systemd/are-they-hiring-web.container
[Unit]
Description=Are They Hiring - Web (FastAPI)
Requires=are-they-hiring-db.service are-they-hiring-ollama.service
After=are-they-hiring-db.service are-they-hiring-ollama.service

[Container]
Image=localhost/are-they-hiring-web:latest
Pod=are-they-hiring.pod
EnvironmentFile=%h/.config/are-they-hiring/.env

[Service]
Restart=always

[Install]
WantedBy=default.target
